# Config/etc/profile.d/alias.sh
#!/usr/bin/bash

alias nano="nano -SLlwxiE --tabsize=4"
