# Config/etc/profile.d/alias.sh
#!/usr/bin/bash

export EDITOR="nano"

alias nano="nano -SLlwxiE --tabsize=4"
