# Config/etc/profile.d/history.sh
#!/usr/bin/bash

export HISTSIZE=500
export HISTFILESIZE=0
export HISTFILE="/dev/null"

export LESSKEY="/etc/sysless"
export LESSHISTFILE="/dev/null"

export PYTHONPYCACHEPREFIX="/var/cache/python"
