# Config/etc/systemd/system/archlinux-keyring-wkd-sync.service.d/override.conf
[Service]
ExecStart      =
ExecStart      = -/usr/bin/archlinux-keyring-wkd-sync
ReadWritePaths = /var/db/pacman/gnupg/
