# Config/etc/systemd/system/ldconfig.service.d/override.conf
[Service]
ExecStart =
ExecStart = /sbin/ldconfig -XC /var/cache/ld.so.cache
